# WebAssignment/DeviceScript.cron
0,15,30,45 * * * * /usr/bin/python home/pi/Desktop/ELSpring2018/WebAssignment/TempReader.py